# conf/web.conf
[program:web]
command=/usr/bin/gunicorn wsgi -c /conf/gunicorn.config.py
directory=/home/xiao_web13/yhj
autostart=true
autorestart=true